# cron/crontab
* * * * * cd /var/www/laravel_docker && /usr/local/bin/php artisan schedule:run >> /var/log/cron.log 2>&1
# An empty line is required at the end of this file for a valid cron file.